# systemd/system/openvpn@.service.d/sangoma_openvpn_override.conf
[Service]
Environment="OPENSSL_ENABLE_MD5_VERIFY=1 NSS_HASH_ALG_SUPPORT=+MD5"
