# systemd/system/chronyd.service.d/sangoma_chrony_override.conf
[Service]
ExecStart=
ExecStart=/usr/sbin/chronyd -f /etc/sangoma_chrony.conf $OPTIONS
